# scripts/ctrader-bot@.service
[Unit]
Description=cTrader FIX Trading Bot
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=600
StartLimitBurst=5

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Documents/ctrader_trading_bot
Environment="PATH=/home/%i/Documents/ctrader_trading_bot/.venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/%i/Documents/ctrader_trading_bot/.env

# Start command
ExecStart=/home/%i/Documents/ctrader_trading_bot/.venv/bin/python3 /home/%i/Documents/ctrader_trading_bot/ctrader_ddqn_paper.py

# Restart policy - CRITICAL for production
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits (adjust based on your needs)
MemoryLimit=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ctrader-bot

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/%i/Documents/ctrader_trading_bot

# Health monitoring via systemd
# Will restart if no heartbeat received within 180s
WatchdogSec=180

[Install]
WantedBy=multi-user.target
